# systemd/agent-remote-node.service
[Unit]
Description=agent-remote node worker
After=network-online.target agent-remote-runtime.service
Wants=network-online.target
Wants=agent-remote-runtime.service

[Service]
Type=simple
EnvironmentFile=-/etc/agent-remote-node/agent-remote-node.env
ExecStart=/usr/local/bin/agent-remote-node run --config /etc/agent-remote-node/config.json
Restart=always
RestartSec=5
User=agent-remote
Group=agent-remote
WorkingDirectory=/var/lib/agent-remote-node
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/agent-remote-node /var/lib/agent-remote /run/agent-remote
ProtectHome=true

[Install]
WantedBy=multi-user.target
